# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:pilot-server]
command=/App/Ascon.Pilot.Daemon /mnt/vol1/settings.xml
