# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:pilot-server]
command=/usr/bin/dotnet /App/Ascon.Pilot.Daemon.dll /mnt/vol1/settings.xml
